# nic-debian/DEBIAN/prerm
#!/bin/bash
systemctl stop amd-nic-metrics-exporter.service

systemctl disable amd-nic-metrics-exporter.service
